# roles/desktop/templates/clean_downloads.timer
[Unit]
Description=Clean downloads timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
