# roles/desktop/templates/restart_pcscd.service
[Unit]
Description=Restart pcscd service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart pcscd

[Install]
WantedBy=default.target
